# SOURCES/deaddrive.upstart
# dead drive

start on stopped rc RUNLEVEL=[2345]
stop on starting rc RUNLEVEL=[S016]

respawn 
exec /usr/local/sbin/deaddrive
